# systemd/mailcheck.service
[Unit]
Description=Mail Checker Daemon 

[Service]
ExecStart=python3 -u 'Path to mailchecker'
Restart=always

[Install]
WantedBy=default.target
